# door.service
# this is a backup of the original service for reference
[Unit]
Description=Door Control Service
After=network.target

[Service]
ExecStart=/usr/bin/python3 /home/pi/start.py
WorkingDirectory=/home/pi/
Restart=always
User=pi
Environment=PYTHONUNBUFFERED=1

[Install]
WantedBy=multi-user.target
